# infra/supervisor-grpc.ini
[program:grpc]
command = /home/app/entrypoint/run_grpc.py --workers=8
user = web
killasgroup = true

stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
