# infra/supervisor-worker.ini
[program:celery]
command = celery --app alvinchow_backend.app.worker worker -l info -c4  --without-mingle --without-gossip --without-heartbeat --max-tasks-per-child=10000
user = nobody
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
killasgroup = true
stopasgroup = true
